# templates/app.sh.erb
#!/bin/bash
# 
# chkconfig: <%= @start_levels %> <%= @start_priority %> <%= @end_priority %>
# description: <%= @description %>

# Import function library
. /etc/rc.d/init.d/functions

prog=`basename $0`

pid_file="<%= @pid_file %>"
stdout_log="<%= @stdout_log %>"
stderr_log="<%= @stderr_log %>"

# Commands
user="<%= @user %>"
proc="<%= @proc %>"
options="<%= @options %>"

case "$1" in
  start)
    echo -n "Starting $prog: "
    daemon --user=$user $proc $options
    RETVAL=$?
    pidofproc $proc > $pid_file
    echo
    exit $RETVAL
  ;;
  stop)
    echo -n "Stopping $prog: "
    killproc -p $pid_file
    RETVAL=$?
    echo
    exit $RETVAL
    ;;
  restart)
    $0 stop && $0 start
    ;;
  status)
    status -p $pid_file $prog
    RETVAL=$?
    exit $RETVAL
    ;;
  *)
    echo "Usage: $0 {start|stop|restart|status}"
    exit 1
    ;;
esac

exit 0
